# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/run/supervisord.pid

[program:php-fpm]
command=php-fpm -F
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20

; ----------------------------------------------------------------------------
; queue-worker — async WhatsApp delivery pipeline
; ----------------------------------------------------------------------------
; Drains the `high` and `default` queues from the database driver (jobs
; table). Single concurrency by design — the old approve() template loop had
; an explicit sleep(1) between Meta calls to dodge rate limits; one worker
; processing one job at a time replicates that pacing with zero in-request
; wait for the admin.
;
;   --tries=3 / --backoff=5,15,60 : 3 attempts, then the row moves to
;                                   failed_jobs for manual replay.
;   --max-time=3600               : worker self-restarts hourly so a deploy
;                                   can't strand it with stale bytecode.
[program:queue-worker]
command=php /app/artisan queue:work --tries=3 --backoff=5,15,60 --max-time=3600 --sleep=3 --queue=high,default
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=30
